# anacron/files/anacrontab
# /opt/etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=/opt/sbin:/opt/bin:/sbin:/bin:/usr/sbin:/usr/bin

# These replace cron's entries
1	5	cron.daily	/opt/sbin/run-anacron /opt/etc/cron.daily
7	10	cron.weekly	/opt/sbin/run-anacron /opt/etc/cron.weekly
@monthly	15	cron.monthly	/opt/sbin/run-anacron /opt/etc/cron.monthly
